# deploy/systemd/xyz-mysql.service
[Unit]
Description=XYZ Agent MySQL (Docker)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=__DEPLOY_USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/docker compose up -d mysql
ExecStop=/usr/bin/docker compose down
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
