# deploy/systemd/xyz-agent-api.service
[Unit]
Description=XYZ Agent FastAPI Service
After=network.target xyz-mysql.service xyz-agent-mcp.service
Requires=xyz-agent-mcp.service
Wants=xyz-mysql.service

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__PROJECT_DIR__
Environment="PATH=__PROJECT_DIR__/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=__PROJECT_DIR__/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
