# deploy/systemd/xyz-agent-poller.service
[Unit]
Description=XYZ Agent Module Poller Service
After=network.target xyz-agent-api.service
Requires=xyz-agent-api.service

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__PROJECT_DIR__
Environment="PATH=__PROJECT_DIR__/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=__PROJECT_DIR__/.venv/bin/python -m xyz_agent_context.services.module_poller
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
