# deploy/systemd/xyz-agent-mcp.service
[Unit]
Description=XYZ Agent MCP Servers
After=network.target

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__PROJECT_DIR__
Environment="PATH=__PROJECT_DIR__/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=__PROJECT_DIR__/.venv/bin/python src/xyz_agent_context/module/module_runner.py mcp
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
